# supervisord.conf
[supervisord]
nodaemon=true
[program:log-sync]
process_name=%(program_name)s_%(process_num)02d
command=php /app/sync.php
redirect_stderr=true
autostart=true
autorestart=true
user=www-data
numprocs=1
